# scripts/com.warm-outreach.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.warm-outreach.morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ari/Documents/zeroclaw-outreach/scripts/run-warm-outreach.sh</string>
    </array>

    <!-- Daily at 7:30 PM (local time — set TZ below for Athens) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Europe/Athens</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ari</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ari/Library/Logs/warm-outreach.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ari/Library/Logs/warm-outreach.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
